# deploy/mywebapp.socket
[Unit]
Description=mywebapp socket

[Socket]
ListenStream=127.0.0.1:3000
NoDelay=true
Service=mywebapp.service

[Install]
WantedBy=sockets.target
